# deploy/systemd/termino-watchdog.service
# PARTIALLY RECONSTRUCTED: the watchdog SCRIPT (path/behaviour/exit-codes) is
# verified against CT 131, but the live .service unit itself was not captured.
# Verify User=/ExecStart against the running unit:
#     systemctl cat termino-watchdog.service
# The live script is root-owned at /usr/local/bin/termino_watchdog.sh, so this
# leaves User unset (runs as root); adjust if your unit differs.
[Unit]
Description=Termino watchdog - alert if today's daily run did not succeed
After=network-online.target

[Service]
Type=oneshot
# ntfy.sh topic for alerts (half-secret). The script also accepts NTFY_TOPIC
# from its own default if this is unset.
Environment="NTFY_TOPIC=<NTFY_TOPIC>"
ExecStart=/usr/local/bin/termino_watchdog.sh
StandardOutput=journal
StandardError=journal
