# deploy/systemd/termino-watchdog.timer
# RECONSTRUCTED from the behaviour spec (runs 07:00 Europe/Vienna, one hour
# after the daily run). Verify against the live unit:
#     systemctl cat termino-watchdog.timer
[Unit]
Description=Run the Termino watchdog daily at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00 Europe/Vienna
Persistent=true

[Install]
WantedBy=timers.target
